# build/Lede_source/files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.
cat /etc/sihp1020.dl >>/dev/usb/lp0
/etc/hotplug.d/usb/10_usb-printer
exit 0
